# worker/celery_worker.conf
[program:celery]
command=celery worker -A celery-runner --loglevel=info
directory=/worker
